# udev_rules/hdmi.rules
ACTION="change", SUBSYSTEM="drm", HOTPLUG="1", RUN+="/usr/local/bin/custom_scripts/hdmi_autoconnect.sh"
